# etc/cron.d/hello-cron
#*/15 * * * * root date >> /var/log/cron.log 2>&1
#*/15 * * * * root /usr/bin/php /var/www/html/cron/status.cron.php
* * * * * root date >> /var/log/cron.log 2>&1
* * * * * root /usr/bin/php /var/www/html/cron/status.cron.php
# An empty line is required at the end of this file for a valid cron file.